# fundflow.service
[Unit]
Description=FundFlow Telegram Bot
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/.gemini/antigravity/scratch/fundflow
Environment="DATABASE_URL=postgresql://fundflow:fundflow@localhost:5432/fundflow"
EnvironmentFile=/root/.gemini/antigravity/scratch/fundflow/config/.env
ExecStart=/root/.gemini/antigravity/scratch/fundflow/venv/bin/python bot/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
